# build/kbsearch_2.0.0_all/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create system user
        if ! id kbsearch >/dev/null 2>&1; then
            useradd -r -s /bin/false -d /var/lib/kbsearch -c "KB Search Service" kbsearch
        fi
        
        # Create directories
        mkdir -p /etc/kbsearch/{run,logs} /var/{lib,log}/kbsearch /var/backups/kbsearch /mnt/kb/{_Gare,_AQ}
        chown -R kbsearch:kbsearch /var/lib/kbsearch /var/log/kbsearch /var/backups/kbsearch /etc/kbsearch/{run,logs}
        
        # Set permissions
        chmod 755 /usr/bin/kb-*
        
        # Generate config if needed
        if [[ ! -f /etc/kbsearch/.env ]]; then
            cp /etc/kbsearch/.env.example /etc/kbsearch/.env
            ADMIN_TOKEN=$(openssl rand -hex 32)
            sed -i "s/ADMIN_TOKEN=.*/ADMIN_TOKEN=$ADMIN_TOKEN/" /etc/kbsearch/.env
            chown kbsearch:kbsearch /etc/kbsearch/.env
            chmod 600 /etc/kbsearch/.env
            echo "Admin token: $ADMIN_TOKEN" > /var/log/kbsearch/install.log
            chmod 600 /var/log/kbsearch/install.log
        fi
        
        systemctl daemon-reload
        systemctl enable kbsearch.service
        
        echo ""
        echo "🎉 KB Search installed successfully!"
        echo "Next steps:"
        echo "  1. Configure: /etc/kbsearch/.env"
        echo "  2. Start: systemctl start kbsearch"
        echo "  3. Status: kb-setup status"
        ;;
esac
exit 0
